# recipes-devtools/python/files/mqtt-db.service
[Unit]
Description=Python MQTT and SQLite3 Application

[Service]
ExecStart=/opt/mqttApp/app.py

[Install]
WantedBy=multi-user.target
